# achifarm.service
[Unit]
Description=Achi Monitoring Service

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
Restart=always

